# quadlets/ollama.container.in
[Unit]
Description=Ollama Application Container

[Container]
Image=@OLLAMA_IMAGE@
AutoUpdate=registry
ContainerName=ollama
Pod=podman-ai-stack.pod
Label=app=podman-ai-stack
Volume=ollama:/root/.ollama:Z
AddDevice=nvidia.com/gpu=all
DropCapability=all
ReadOnly=true
Environment=OLLAMA_HOST=0.0.0.0
EnvironmentFile=/etc/sysconfig/podman-ai-stack
#EnvironmentFile=-%h/.config/podman-ai-stack.env
PodmanArgs=--memory=@OLLAMA_MEMORY@ --cpus=@OLLAMA_CPUS@
HealthCmd=ollama list || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=90s

[Service]
Restart=on-failure
ExecStartPre=/usr/bin/podman pull @OLLAMA_IMAGE@

[Install]
WantedBy=default.target