# quadlets/open-webui.container.in
[Unit]
Description=Open WebUI Application Container

[Container]
Image=@OPEN_WEBUI_IMAGE@
AutoUpdate=registry
ContainerName=open-webui
Pod=podman-ai-stack.pod
Label=app=podman-ai-stack
Environment=OLLAMA_BASE_URL=http://localhost:11434
EnvironmentFile=/etc/sysconfig/podman-ai-stack
#EnvironmentFile=-%h/.config/podman-ai-stack.env
DropCapability=all
ReadOnly=true
# Securely pass Database URLs or API keys via Podman Secrets.
# Uncomment the lines below after creating the secrets via 'podman secret create':
# Secret=openwebui_database_url,type=env,target=DATABASE_URL
# Secret=openai_api_key,type=env,target=OPENAI_API_KEY
Volume=open-webui:/app/backend/data:Z
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:8080 || exit 1
HealthStartPeriod=60s
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=on-failure
ExecStartPre=/usr/bin/podman pull @OPEN_WEBUI_IMAGE@

[Install]
WantedBy=default.target
